# docker/crontab
# run on container startup
@reboot /usr/local/bin/porkcron
@weekly /usr/local/bin/porkcron
