# systemd/porkcron.timer
[Unit]
Description=Run porkcron weekly

[Timer]
OnCalendar=Weekly
Persistent=true

[Install]
WantedBy=timers.target
